# src/ofono.service.in
[Unit]
Description=Telephony service
After=syslog.target

[Service]
Type=dbus
BusName=org.ofono
User=radio
EnvironmentFile=-/var/lib/environment/ofono/*.conf
ExecStart=@prefix@/sbin/ofonod -n $OFONO_ARGS
StandardError=null
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
